# Build/plist/yanue.v2rayu.sing-box.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>yanue.v2rayu.sing-box</string>
	<key>ProgramArguments</key>
	<array>
		<string>__SINGBOX_BIN__</string>
		<string>run</string>
		<string>-c</string>
		<string>__APP_HOME_DIR__/config.json</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>KeepAlive</key>
	<true/>
	<key>WorkingDirectory</key>
	<string>__APP_HOME_DIR__</string>
	<key>StandardOutPath</key>
	<string>__APP_HOME_DIR__/run-sing-box.log</string>
	<key>StandardErrorPath</key>
	<string>__APP_HOME_DIR__/run-sing-box.log</string>
</dict>
</plist>
